# crab.service
[Unit]
Description=Crab daemon

[Service]
ExecStart=/usr/bin/crab
Restart=on-failure

[Install]
WantedBy=multi-user.target 